# install/systemd/matts-console.service
[Unit]
Description=MDE LLM-PROXY V2 React Console
Documentation=https://github.com/matthewmackes/MACKESCODE-AI-PROXY
After=network.target matts-value-set-proxy.service
Wants=network.target matts-value-set-proxy.service

[Service]
Type=simple
User=matts
Group=matts
WorkingDirectory=/var/lib/matts-value-set
EnvironmentFile=/etc/matts-value-set/environment.conf
Environment=PYTHONPATH=/usr/lib/matts-value-set/vendor/python:/usr/lib/matts-value-set
Environment=TMUX_TMPDIR=/var/lib/matts-value-set/tmux
ExecStart=/usr/bin/python3 /usr/lib/matts-value-set/matts-v2-console.py \
  --host 0.0.0.0 \
  --port 18182
Restart=always
RestartSec=10
StartLimitInterval=60
StartLimitBurst=5

# Security hardening
NoNewPrivileges=true
ProtectSystem=strict
ProtectHome=true
ReadWritePaths=/var/lib/matts-value-set /var/log/matts-value-set /tmp
PrivateTmp=true
PrivateDevices=true
ProtectKernelTunables=true
ProtectKernelModules=true
ProtectControlGroups=true

# Resource limits
LimitNOFILE=65536
LimitNPROC=4096

# Logging
StandardOutput=journal
StandardError=journal
SyslogIdentifier=matts-v2-console

[Install]
WantedBy=multi-user.target
